# unrestricted-recovery/accelerated-container-startup.path
[Path]
PathExists=/var/lib/kubelet/cpu_manager_state

[Install]
WantedBy=multi-user.target
